# deploy/systemd/dns-control-api.service
[Unit]
Description=DNS Control API v2.1
After=network.target
Documentation=https://github.com/ocaiobarros/dns-orchestrator
ConditionPathIsDirectory=/etc/nftables.d
ConditionPathExists=/opt/dns-control/backend/venv/bin/uvicorn

[Service]
Type=simple
User=dns-control
Group=dns-control
WorkingDirectory=/opt/dns-control/backend
EnvironmentFile=/etc/dns-control/env
ExecStartPre=/usr/bin/test -d /etc/nftables.d
ExecStartPre=/usr/bin/test -x /opt/dns-control/backend/venv/bin/uvicorn
ExecStart=/opt/dns-control/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dns-control-api

# Security hardening
# NOTE: All paths in ReadWritePaths MUST be created by the installer
# before this unit starts, or ProtectSystem=strict causes status=226/NAMESPACE.
# The installer (install_debian13.sh) creates these in Step 3.
NoNewPrivileges=false
ProtectSystem=strict
ReadWritePaths=/var/lib/dns-control /var/log/dns-control /etc/unbound /etc/nftables.d /etc/nftables.conf /etc/network /etc/systemd/system /etc/sysctl.d /etc/frr /etc/default /usr/lib/systemd/system
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
